# rustinx.service
[Unit]
Description=Rustinx for monitoring Nginx

[Service]
ExecStart=/usr/bin/rustinx
Restart=always

[Install]
WantedBy=multi-user.target
